# root/etc/crontabs/abc
# min                                hour    day     month   weekday command
0                                     */6     *       *       *       /usr/bin/php /app/www/public/crons/cleanup.php
30                                    23      *       *       *       /usr/bin/php /app/www/public/crons/backup.php